# infra/files/wispgen-backup.service
[Unit]
Description=Backup WISPGen tenant databases to S3

[Service]
Type=oneshot
User=wispgen
Group=wispgen
ExecStart=/srv/wispgen/src/infra/files/wispgen-backup.sh
